# aafai-bus.service
[Unit]
Description=Gunicorn instance to serve aafai-bus
After=network.target

[Service]
# Environment variables for the application
Environment="APP_ENV=production" "QUEUE_BASE_PATH=/var/www/aafai-bus/data"

# Run the service as the 'www-data' user
User=www-data
Group=www-data

# Set the working directory to your app's root
WorkingDirectory=/var/www/aafai-bus

# Command to start the app using Gunicorn and the application factory
# The 'server:create_app()' part means "in server.py, call the function create_app()"
ExecStart=/var/www/aafai-bus/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 "server:create_app()"

# Restart the service if it ever fails
Restart=always

[Install]
WantedBy=multi-user.target
